# launchd/user-agents/com.northridge.hyperdeck.python-sync-smoke.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.northridge.hyperdeck.python-sync-smoke</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/northridge/scripts/launchd/run_python_sync_smoke_test.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/northridge/scripts</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/northridge</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/northridge/Library/Logs/hyperdeck/python-launchd-sync.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/northridge/Library/Logs/hyperdeck/python-launchd-sync.log</string>

  <key>LimitLoadToSessionType</key>
  <array>
    <string>Aqua</string>
  </array>
</dict>
</plist>
